# etc/cron.d/eos-health
*/5 * * * * root /etc/init.d/eoshealth status >& /dev/null




