# etc/cron.d/xrd-logs
30 0 * * * root  find /var/log/xroot  -name "xrdlog.*.????????" -exec gzip -9 \{\} \; ; find /var/log/xroot  -name "xrdlog.*.????????.gz" -mtime +90 -exec rm \{\} \;

