# misc/etc/systemd/system/eos.target
[Unit]
Description=eos target allowing to start/stop all eos@*.service instances at once
[Install]
WantedBy=multi-user.target
