# misc/etc/systemd/system/eos@.socket
[Unit]
Description=Eos daemon socket

[Sockets]
# Set the xrd port
ListenStream=1094
# Set the http port
ListenStream=900
Service=eos@%i.service
